# collector/c5500xk-collector.service
[Unit]
Description=C5500XK direct Bluetooth collector
After=bluetooth.service network-online.target
Wants=bluetooth.service network-online.target

[Service]
Type=simple
ExecStart=/opt/c5500xk-collector/venv/bin/python /opt/c5500xk-collector/c5500xk_collector.py --config /etc/c5500xk-collector/config.json
Restart=on-failure
RestartSec=10
User=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/etc/c5500xk-collector/config.json

[Install]
WantedBy=multi-user.target
